# rootfs_5/ufw/rootfs/usr/lib/opkg/info/hpsrv.list
/etc/init.d/hpsrv
/sbin/hpsrv
/etc/config/hpsrv
