# rootfs_5/ufw/rootfs/usr/lib/opkg/info/log.list
/etc/config/log
/sbin/log
/etc/init.d/log
