# rootfs_5/ufw/rootfs/usr/lib/opkg/info/proxy.list
/etc/config/proxy
/etc/init.d/proxy
/sbin/proxy
