# rootfs_5/ufw/rootfs/usr/lib/opkg/info/uappd.list
/etc/init.d/uappd
/sbin/uappd
/etc/config/uappd
