# rootfs_3/saf/rootfs/usr/lib/opkg/info/cloudclient.list
/etc/init.d/cloudclt
/etc/config/cloudclient
/rom/dbus_map
/sbin/cloudclient
/sbin/cloudclocal
