# rootfs_3/saf/rootfs/usr/lib/opkg/info/dsmonitor.list
/etc/init.d/dsmonitor
/usr/sbin/dsmonitor
