# rootfs_3/saf/rootfs/usr/lib/opkg/info/easymesh.list
/sbin/easymesh
/etc/init.d/easymesh
/etc/config/easymesh
